# manifest/etc/supervisor/conf.d/exporter.conf
[program:exporter]
command=/usr/bin/dmarcts-report-exporter.pl
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
